# dropbox/.config/systemd/user/dropbox.service
[Unit]
Description=Dropbox Headless

[Service]
ExecStart=/home/tie/.dropbox-dist/dropboxd
Restart=on-failure

[Install]
WantedBy=default.target
